# monitor.cron
MAILTO=me@yianni.dev
0 * * * * /usr/local/bin/change.py
