# monitor.cron
MAILTO=me@yianni.dev
* * * * * /usr/local/bin/change.pl
